# service/phpfpm/run
#!/usr/bin/env bash
exec /usr/local/sbin/php-fpm -c /usr/local/etc/php/ -y /usr/local/etc/php-fpm.conf